# cronjob/update_task_status.cron
* * * * * /usr/local/bin/python /app/manage.py update_task_status